# agva-atv-safari/output/results/supervisord.conf
[supervisord]
nodaemon=false
logfile=/home/user/webapp/agva-atv-safari/output/results/supervisord.log
pidfile=/home/user/webapp/agva-atv-safari/output/results/supervisord.pid

[unix_http_server]
file=/home/user/webapp/agva-atv-safari/output/results/supervisor.sock

[supervisorctl]
serverurl=unix:///home/user/webapp/agva-atv-safari/output/results/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:webserver]
command=python3 /home/user/webapp/agva-atv-safari/output/results/web_server.py
directory=/home/user/webapp/agva-atv-safari/output/results
autostart=true
autorestart=true
stdout_logfile=/home/user/webapp/agva-atv-safari/output/results/webserver.log
stderr_logfile=/home/user/webapp/agva-atv-safari/output/results/webserver_error.log
